# quadlets/cognicode-mcp.container
[Container]
ContainerName=cognicode-mcp
Image=cognicode-mcp:latest

# Build from local Dockerfile
# Rebuild: just build-image && systemctl --user restart cognicode-mcp
SetImageCommand=/usr/bin/podman build --tag cognicode-mcp:latest /var/home/rubentxu/Proyectos/rust/CogniCode

# Connect to PostgreSQL container via podman network
# The PG container is accessible as 'cognicode-postgres' on the podman bridge
Environment=DATABASE_URL=postgres://cognicode:cognicode@cognicode-postgres:5432/cognicode
Environment=COGNICODE_PROJECT=/workspace
Environment=RUST_LOG=info

# Mount the project workspace (read-only)
Volume=%h/Proyectos:/workspace:ro

# Expose MCP via HTTP/SSE on port 9847
# OpenCode connects as remote MCP: http://localhost:9847/mcp
PublishPort=127.0.0.1:9847:9847

# Depends on PostgreSQL being healthy
DependsOn=cognicode-postgres.service

# Auto-restart
[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
